# Systems/eol-rt-data/ensure_running.cron
#!/bin/csh

ps ax | grep udp2sql | grep -qv grep
if ($status) then                                                                                                                                                                                    
  /opt/local/bin/udp2sql >& /tmp/udp2sql.log &
endif
