# Systems/eol-rt-data/ensure_running.cron
#!/bin/csh

ps ax | grep udp2sql | grep -qv grep
if ($status) then                                                                                                                                                                                    
  /home/local/Systems/eol-rt-data/udp2sql/src/run_udp2sql.sh
endif
