# Systems/eol-rt-data/ensure_running.cron
#!/bin/csh

ps ax | grep udp2sql | grep -qv grep
if ($status) then                                                                                                                                                                                    
  nohup /home/ads/udp2sql/src/udp2sql >& /var/log/udp2sql.log &
endif
