# Systems/crontabs/crontab.ads.tikal
#
# NWS Radar for ground OpenStreetMaps.  Adjust frequency on a per project basis.
# Maybe 15 minutes when no projects, and 3-5 minutes when US project in progress.
*/10 * * * * /h/eol/ads/crontab/osm_images.cron

