# Systems/crontabs/crontab.ads.gv
SHELL=/bin/csh

@reboot /home/local/Systems/scripts/scrubber >& /tmp/scrubber.log
@reboot ssh -N eoltunnel >& /tmp/eoltunnel.log

# Always
*/3 * * * * /home/local/Systems/scripts/ensure_running.cron >& /tmp/ensure_run.log
*/2 * * * * /home/local/Systems/scripts/send_camera.cron >& /tmp/camera_send.log
0 * * * * /home/local/Systems/scripts/scrubber >& /tmp/scrubber.log

# New satellite script
3,18,33,48 * * * * /home/local/Systems/scripts/get_ir_image_cron.py >>& /tmp/get_ir_image.log 
5,20,35,50 * * * * /home/local/Systems/scripts/get_4kmvis_image_cron.py >>& /tmp/get_4kmvis_image.log 
#5,20,35,50 * * * * /home/local/Systems/scripts/get_vis_image_cron.py >>& /tmp/get_vis_image.log 
*/10 * * * * /home/local/Systems/scripts/get_conus_cweb_cron.py >>& /tmp/get_conus_cweb.log

# Get flight plan script
#*/20 * * * * /home/local/Systems/scripts/get_flight_plan.py >>& /tmp/get_flight_plan.log 
#*/30 * * * * /home/local/Systems/scripts/get_drop_plan.py >>& /tmp/get_drop_plan.log 

# MTP
#*/3 * * * * /home/local/Systems/scripts/send_MTP.cron >& /tmp/mtp_send.log 
# TOGA
*/2 * * * * python /home/local/Systems/scripts/send_toga.py >>& /tmp/send_toga.log
# AVAPS
#* * * * * /home/local/Systems/scripts/avaps_proc.cron.py >>& /tmp/avaps_proc.log 

# 
# Persistent ssh tunnel to EOL, attempting reopens only every 10 minutes.
*/10 * * * * ssh -N eoltunnel >& /tmp/eoltunnel.log

# Time stamp label for MC
* * * * * /home/local/Systems/scripts/make_label.py >>& /dev/null

# Copy flight_track and position to CatalogIngest
* * * * * /home/local/Systems/scripts/copy_flight_geo_to_catalog.cron >>& /tmp/copy_flight_geo_to_catalog.log

# Cleanup old Catalog files from filesystem and database via DockerCompose
1 */2 * * * CATALOG_MAPS_FILE_DELETE_CONFIRMATION_HASH=6d854e052d3f061918af6fa8ed0e0f209b6893793813aaac4dd952d8bb1ed9a6 sudo /usr/local/bin/docker-compose -f /home/catalog/catalog-maps/docker-compose.yml run -e CATALOG_MAPS_FILE_DELETE_CONFIRMATION_HASH=$CATALOG_MAPS_FILE_DELETE_CONFIRMATION_HASH maps ./bin/rails runner etc/acserver/delete-catalog-files.rb

# IRC and IRC bot
@reboot /home/ads/ircd/Unreal3.2/unreal start >& /dev/null
@reboot sleep 5; /home/ads/rbot/bin/rbot >& /dev/null
