# Systems/crontabs/crontab.root.ac
*/5 * * * * /usr/sbin/ddclient -verbose -file /etc/ddclient/c130.conf >> /tmp/ddclient.log 2>&1
