# Systems/Systems/crontabs/crontab.ads.barolo
*/5 * * * * /h/eol/ads/crontab/ensure_running > /tmp/ads.crontab.log
0 5 * * mon /net/jlocal/raf/Systems/scripts/backupCalib.cron >& /dev/null
#
#
# This is run on tikal.  NWS Radar for ground OpenStreetMaps
#*/5 * * * * /h/eol/ads/crontab/osm_images.cron
