# Systems/scripts/catalog_delete_files.cron
#!/bin/bash

#
# Cleanup old Catalog files from filesystem and database via DockerCompose
#

CATALOG_MAPS_FILE_DELETE_CONFIRMATION_HASH=6d854e052d3f061918af6fa8ed0e0f209b6893793813aaac4dd952d8bb1ed9a6

/usr/local/bin/docker-compose -f /home/catalog/catalog-maps/docker-compose.yml run -e CATALOG_MAPS_FILE_DELETE_CONFIRMATION_HASH=$CATALOG_MAPS_FILE_DELETE_CONFIRMATION_HASH maps ./bin/rails runner etc/acserver/delete-catalog-files.rb
