# Systems/scripts/ensure_running.cron
#!/bin/bash

# Script to ensure miscellanious prcesses are running.  Run out of cron
# periodically, say every 3-5 minutes.

if ! ps ax | grep status_listener | grep -qv grep
then
  nohup status_listener 2>&1 /tmp/status_listener.log &
fi

if ! ps ax | grep acTrack2kml | grep -qv grep
then
  nohup acTrack2kml -o 2>&1 /tmp/acTrack2kml.log &
fi

if ! ps ax | grep listen_sat.pl | grep -qv grep
then
  nohup listen_sat.pl 2>&1 /tmp/listen_sat.log &
fi
