# Systems/scripts/get_wx.cron
#!/bin/bash

cd /var/www/html

pkill wget

# Limit rate to 30kbps... assuming mpds link speed varies between 20-50kbps.
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_Vis.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_WV.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_IR.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg1_Vis.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg1_WV.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg1_IR.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg2_Vis.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg2_WV.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg2_IR.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg3_Vis.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg3_WV.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg3_IR.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg45_Vis.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg45_WV.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg45_IR.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg6_Vis.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg6_WV.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg6_IR.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg7_Vis.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg7_WV.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg7_IR.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg8_Vis.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg8_WV.jpg
wget -N --limit-rate=3k http://catalog.eol.ucar.edu/hippo/other/latest_goes-11_leg8_IR.jpg
