# Configuration/crontabs/crontab.ads.steam
*/3 * * * * /home/local/projects/scripts/ensure_running_lab.cron >/dev/null 2>&1

#*/10 * * * * ssh -N eoltunnel > /tmp/eoltunnel.log 2>&1
