# Configuration/crontabs/crontab.ads.eol-rt-data
* * * * * /home/local/eol-rt-data/scripts/ensure_running.cron >/dev/null 2>&1
